# systemd/brokerai-update.service
[Unit]
Description=BrokerAI auto-update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/brokerai/config.env
ExecStart=/opt/brokerai/scripts/auto-update.sh
