# systemd/brokerai-postgres-backup.timer
[Unit]
Description=BrokerAI daily Postgres backup

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
